# campus.conf
[program:campus]
command = /home/jiaqing/webapps/admin_campus/gunicorn.sh
stdout_logfile = /home/jiaqing/webapps/admin_campus/logs/gunicorn_supervisor.log
redirect_stderr = true
