# seed/wiki.conf
[program:sub_domain]
command = server_root/sub_domain/gunicorn_wiki.sh
stdout_logfile = server_root/sub_domain/logs/gunicorn_supervisor.log
redirect_stderr = true
